# com.finance.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">

<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.finance.dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/saatvikkamjula/finance-dashboard/run_dashboard_scheduled.command</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/saatvikkamjula/finance-dashboard</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/saatvikkamjula/finance-dashboard/reports/scheduler_output.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/saatvikkamjula/finance-dashboard/reports/scheduler_error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
